# src/setup/utility/scripts/CDRPush.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

exec su - <INSTALL USER> -c "<INSTALLROOT>/cdrSftpUtility/dbPush/scripts/startDbPush.sh <INSTALLROOT>/cdrSftpUtility/dbPush"

