# src/setup/utility/scripts/CDRSftp.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn

exec su - <INSTALL USER> -c "<INSTALLROOT>/cdrSftpUtility/fileSftp/scripts/startFileSftp.sh <INSTALLROOT>/cdrSftpUtility/fileSftp"

